# misc/tuxedo-tomted.service
[Unit]
Description=Tomte-daemon, finishes tasks that could not be accomplished

[Service]
ExecStart=/usr/bin/tuxedo-tomte configure all
Restart=no
