# misc/tuxedo-tomted.timer
[Unit]
Description: Starts tuxedo-tomte after waiting for predefined time

[Timer]
OnActiveSec=360
